# pidash.service
[Unit]
Description=PiDash
After=default.target

[Service]
ExecStart=/home/pi/Documents/pidash/env/bin/python /home/pi/Documents/pidash/dashSocket.py
User=pi
Restart=always
Environment="DISPLAY=:0"

[Install]
WantedBy=default.target